# supervisor_scripts/amm_module_manager.conf
[program:amm_module_manager]
autostart=false
autorestart=true
directory=/usr/local/bin
command=/usr/local/bin/amm_module_manager -s -a
priority=105
numprocs=1
user=amm
redirect_stderr=true
stdout_logfile=/usr/local/bin/logs/amm_module_manager.log

